# utils/dump_ryzenadj_table.service
[Unit]
Description=test job

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash ryzenadj -i >> /tmp/ryzenadj_dump
